# files/etc/profile.d/alias.sh
alias fd='find . -type d -name'
alias fdi='find . -type d -iname'
alias ff='find . -type f -name'
alias ffi='find . -type f -iname'
alias l='ls -lFh'
alias lS='ls -1FSsh'
alias la='ls -lAFh'
alias lart='ls -1Fcart'
alias dot='ls -ld .*'
alias ll='ls -l'
alias lr='ls -tRFh'
alias lrt='ls -1Fcrt'
alias lsa='ls -lah'
alias lt='ls -ltFh'
alias md='mkdir -p'
alias mkdir='mkdir -pv'
alias rm='rm -i'
alias ramfree='sync && echo 3 > /proc/sys/vm/drop_caches'
alias shutdown='poweroff'